# upstart_router_watch.conf
start on net-device-up
#start on startup
respawn
respawn limit unlimited

chdir /home/pi/node-raspi-security

script
  node router_reboot.js 2>&1 >> /var/log/router_reboot.log
end script

